# klimbxd.service
[Unit]
Description = My Custom Service

[Service]
Type = forking
ExecStart = /bin/klimbxd
Restart = on-failure

[Install]